# supervisor/pybot.example.conf
[program:pybot]

command=/PATH/TO/pybotd
directory=/PATH/TO/REPOSITORY
user=UNPRIVILEGED_USER
user=UNPRIVILEGED_GROUP

# Do not edit below this line
#############################
numprocs=1
autorestart=true
autostart=true
startsecs=3
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
stopasgroup=false
killasgroup=false
stdout_logfile=/tmp/pybot.stdout.log
stderr_logfile=/tmp/pybot.stdout.log
